# containers/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck
foreground { chmod o+w /dev/stdout }
foreground { chmod o+w /dev/stderr }
s6-setuidgid www-data
/usr/sbin/nginx